# rootfs/etc/services.d/caddy/run
#!/usr/bin/with-contenv sh

env HOME=/srv s6-setuidgid ${TC_USER} /usr/bin/caddy run --config /srv/caddy/config.json
