# rootfs/etc/services.d/trojan/run
#!/usr/bin/with-contenv sh

sleep 5

CERTFILE="/srv/trojan/trojan.crt"

if [ -f $CERTFILE ]; then
    s6-setuidgid ${TC_USER} /usr/local/bin/trojan /srv/trojan/config.json
fi
